# INTANG.service

[Unit]
Description=INTANG
After=multi-user.target

[Service]
Type=forking
ExecStart=/opt/INTANG/bin/intangd 10
KillMode=control-group

[Install]
WantedBy=multi-user.target

